# wifi/deploy/linux/systemd/helen-rendezvous.service
[Unit]
Description=Helen NAT-traversal Rendezvous Coordinator
Documentation=Internal LAN deployment — see DEPLOY-GUIDE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=helen
Group=helen
WorkingDirectory=/opt/helen-rendezvous
ExecStart=/opt/helen-rendezvous/Helen-Rendezvous
EnvironmentFile=/opt/helen-rendezvous/.env
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

LimitNOFILE=65536
LimitNPROC=4096

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/helen-rendezvous/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=helen-rendezvous

[Install]
WantedBy=multi-user.target
