# wifi/deploy/linux/systemd/helen-router.service
[Unit]
Description=Helen Router — mandatory LAN entry point
Documentation=Internal LAN deployment — see DEPLOY-GUIDE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=helen
Group=helen
WorkingDirectory=/opt/helen-router
ExecStart=/opt/helen-router/Helen-Router
EnvironmentFile=/opt/helen-router/.env
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

LimitNOFILE=65536
LimitNPROC=4096

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/helen-router/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=helen-router

[Install]
WantedBy=multi-user.target
